# deployment/spam-classifier-task.conf
[program:spam-classifier-task]
command=/home/spam_classifier/venv/bin/celery -A spam_classifier.celery worker --loglevel=ERROR --concurrency=1
directory=/home/spam_classifier
user=spam_classifier
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
