# deployment/spam-classifier.conf
[program:spam-classifier]
command=/home/spam_classifier/venv/bin/gunicorn -b :8000 -w 4 spam_classifier:app
directory=/home/spam_classifier
user=spam_classifier
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0